# etc/runit/sv/earlyoom/run
#!/bin/sh
[ -r conf ] && . ./conf
exec earlyoom -s 100 $EARLYOOM_ARGS -r 3600 --prefer '^firefox$'
