# etc/runit/sv/earlyoom/run
#!/bin/sh
[ -r conf ] && . ./conf
exec earlyoom $EARLYOOM_ARGS -r 3600
